# mac-proxy/launchd/com.local.ccs-cliproxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.ccs-cliproxy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/ccs</string>
    <string>cliproxy</string>
    <string>start</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/admin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <!-- KeepAlive=false: cliproxy сам даемонизируется, не нужен watch от launchd -->
  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/ccs-cliproxy.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ccs-cliproxy.log</string>
</dict>
</plist>
